# config/containers/systemd/groq.container
[Unit]
Description=Groq AI Dictation Container
After=graphical-session.target sound.target

[Container]
ContainerName=groq
Image=localhost/groq-dictate:latest
# 1. Hardware Access (Keyboard)
AddDevice=/dev/input:/dev/input
# 2. Wayland, DBUS, and PipeWire Access
# %U automatically resolves to your user ID (e.g., 1000)
Volume=%t:%t
EnvironmentFile=%h/.secrets
Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
Environment=PULSE_SERVER=unix:%t/pulse/native
# 3. Security Overrides
SecurityLabelDisable=true
# Pass supplementary groups directly to podman
PodmanArgs=--group-add=keep-groups --userns=keep-id

[Install]
WantedBy=default.target
