# config/containers/systemd/stremio.container
[Unit]
Description=Stremio Streaming Server

[Container]
Image=docker.io/stremio/server:latest
PublishPort=11470:11470
Volume=%h/.local/share/stremio-server:/root/.stremio-server:Z
# Join the shared network
Network=homelab.network
# Give it a predictable hostname for Cloudflare to find
ContainerName=stremio-server

[Install]
WantedBy=default.target
