# config/systemd/user/nextcloud-preview.timer
[Unit]
Description=Run Nextcloud Preview Generation every 15 mins

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
