# archive/ai_agent.container
[Unit]
Description=Persistent AI Agent Container
After=network-online.target

[Container]
ContainerName=ai-agent
Image=ai_agent:latest
Volume=/var/home/eric/src/gh/dcoffline/ai-employee/data:/app/data:Z
EnvironmentFile=%h/.secrets
# Use rootless networking
Network=slirp4netns

[Install]
# Start on user login
WantedBy=default.target
